# deployment/algotread-news.service
[Unit]
Description=AlgoTread News Broadcast Service
After=network.target

[Service]
Type=oneshot
User=YOUR_USERNAME
WorkingDirectory=/path/to/algotread
Environment="PATH=/path/to/algotread/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/path/to/algotread/.env
ExecStart=/path/to/algotread/.venv/bin/python -m scripts.broadcast_news --max 10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=algotread-news
